# files/root/etc/s6-overlay/s6-rc.d/laravel-notification/run
#!/command/with-contenv sh

exec php /var/www/artisan notification:worker