# files/root/etc/s6-overlay/s6-rc.d/laravel-routes/up
#!/command/with-contenv sh

php /var/www/artisan route:cache