# files/root/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv sh

/usr/local/sbin/php-fpm -R