# files/root/etc/s6-overlay/s6-rc.d/laravel-config/up
#!/command/with-contenv sh

php /var/www/artisan config:cache