# files/root/etc/s6-overlay/s6-rc.d/laravel-queue/run
#!/command/execlineb -P

php /var/www/artisan queue:work --sleep=3 --tries=3