# files/root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv sh

/usr/sbin/nginx -g 'daemon off;'