# files/root/etc/s6-overlay/s6-rc.d/laravel-migrations/up
#!/command/with-contenv sh

php /var/www/artisan migrate --force