# files/root/etc/s6-overlay/s6-rc.d/laravel-views/up
#!/command/with-contenv sh

php /var/www/artisan view:cache --no-interaction --quiet