# files/root/etc/s6-overlay/s6-rc.d/laravel-clear/up
#!/command/with-contenv sh

php /var/www/artisan cache:clear --no-interaction --quiet