# files/root/etc/s6-overlay/s6-rc.d/laravel-schedule/run
#!/command/with-contenv sh

php /var/www/artisan schedule:run