# files/root/etc/s6-overlay/s6-rc.d/laravel-schedule/run
#!/command/execlineb -P

cd /var/www

php artisan schedule:work