# files/root/etc/s6-overlay/s6-rc.d/laravel-cache/up
#!/command/with-contenv sh

php /var/www/artisan config:cache &&
php /var/www/artisan route:cache &&
php /var/www/artisan view:cache