# files/root/etc/s6-overlay/s6-rc.d/laravel-cache/up
#!/command/with-contenv sh

exec php /var/www/artisan config:cache
exec php /var/www/artisan route:cache
exec php /var/www/artisan view:cache